# presence-agent/com.vibe.presence-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibe.presence-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/seth/.vibe/presence-agent/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VIBE_API_URL</key>
        <string>https://www.slashvibe.dev</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/seth/.vibe/presence-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seth/.vibe/presence-agent.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
